# ansible/plex/templates/plex-data-sync.timer.j2
{{ ansible_managed | comment }}

[Unit]
Description=Plex Backup Metadata to NFS

[Timer]
Persistent=true
#OnCalendar=*-*-* 00:00:00
OnCalendar=hourly

[Install]
WantedBy=timers.target
